# agents/systemd/raw-watcher.service
[Unit]
Description=Overseer raw session watcher
After=network.target

[Service]
Type=simple
WorkingDirectory=%h/dev/personal/overseer/scripts
ExecStart=/bin/bash %h/dev/personal/overseer/scripts/raw-watcher.sh
Restart=on-failure
RestartSec=15

EnvironmentFile=-%h/.secrets/master.env
Environment=VAULT_PATH=%h/vault

[Install]
WantedBy=default.target
