# agents/systemd/yap_processor.path
[Unit]
Description=Watch vault inbox/yap for new yap files

[Path]
PathModified=%h/vault/inbox/yap
Unit=yap_processor.service

[Install]
WantedBy=default.target
