# agents/systemd/yap_processor.service
[Unit]
Description=Brain Yap Processor Agent
After=network.target

[Service]
Type=simple
WorkingDirectory=%h/dev/personal/overseer/agents
ExecStart=/usr/bin/python3 %h/dev/personal/overseer/agents/yap_processor.py
Restart=on-failure
RestartSec=15

EnvironmentFile=-%h/.secrets/master.env
Environment=VAULT_PATH=%h/vault

[Install]
WantedBy=default.target
